# infra/modules/customer-instance/files/agnes-db-backup.service
[Unit]
Description=Agnes daily system.duckdb backup + canary restore-verify

[Service]
Type=oneshot
ExecStart=/usr/local/bin/agnes-db-backup.sh
