# infra/modules/customer-instance/files/agnes-watchdog.service
[Unit]
Description=Agnes watchdog — incident signature checks

[Service]
Type=oneshot
ExecStart=/usr/local/bin/agnes-watchdog.sh
