# infra/modules/customer-instance/files/agnes-db-backup.timer
[Unit]
Description=Daily Agnes system.duckdb backup (03:17 UTC)

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true

[Install]
WantedBy=timers.target
