# scripts/ops/agnes-state-applier.service
[Unit]
Description=Agnes DB backend state applier — apply target compose lifecycle
After=docker.service agnes-state-applier-bootstrap.service
Wants=docker.service
Requires=agnes-state-applier-bootstrap.service

[Service]
Type=oneshot
User=agnes-applier
Group=docker

# The companion agnes-state-applier-bootstrap.service (root, oneshot,
# RemainAfterExit) creates the agnes-applier user + chowns /data/state
# before this unit's User= directive is resolved. See that unit for
# the rationale — systemd validates User= at unit LOAD time, so we
# can't self-bootstrap inside this unit's own ExecStartPre.

ExecStart=/usr/local/bin/agnes-state-applier.sh
