# scripts/ops/agnes-state-applier-bootstrap.service
[Unit]
Description=Bootstrap state for agnes-state-applier — create user, prep /data/state
Before=agnes-state-applier.service agnes-state-applier.timer
After=local-fs.target

# Runs as root (no User=) so we can actually create the agnes-applier
# user. The companion agnes-state-applier.service is ordered After=
# this one (via Before= here) so by the time systemd loads the User=
# directive there, the user definitely exists.
#
# Phase 8.1 follow-up #2: an earlier attempt put the user-creation in
# ExecStartPre=+ of the main unit. systemd validates User= at unit
# LOAD time, not at ExecStartPre run time, so the unit refused to
# even start with "Failed to determine user credentials" before any
# ExecStartPre had a chance to fire. Splitting into a dedicated
# root-running bootstrap unit is the canonical systemd pattern for
# self-bootstrapping user accounts.

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c 'id -u agnes-applier >/dev/null 2>&1 || useradd --system --no-create-home --shell /usr/sbin/nologin --user-group agnes-applier'
ExecStart=/usr/sbin/usermod -aG docker agnes-applier
ExecStart=/bin/mkdir -p /data/state /data/postgres
ExecStart=/bin/chown -R agnes-applier:agnes-applier /data/state
ExecStart=/bin/bash -c 'chown 70:70 /data/postgres && chmod 700 /data/postgres'
# Make /opt/agnes/.env readable by the applier. The script sources
# this file (IMAGE_REPO, IMAGE_TAG, POSTGRES_PASSWORD, …) AND docker
# compose's --env-file path reads it directly (not via shell source).
# We chown to ``agnes-applier:agnes-applier`` (owner = applier user,
# mode 0600) rather than just chgrp + 0640: docker compose's Go file
# loader opens the file with flags that fail for group-only access
# on some libc versions ("permission denied" on the open() syscall
# even when /etc/group lists the user in the file's group). Setting
# the owner makes it unambiguous regardless of how the caller opens
# the file. Bootstrap unit runs as root, has the privileges to flip
# ownership; the file still contains secrets but is owner-only readable
# (mode 0600), so the security posture is unchanged — only the
# identity of the owner moved from root → agnes-applier.
ExecStart=/bin/bash -c 'if [ -f /opt/agnes/.env ]; then chown agnes-applier:agnes-applier /opt/agnes/.env && chmod 0600 /opt/agnes/.env; fi'

[Install]
WantedBy=multi-user.target
