# services/corporate_memory/systemd/corporate-memory.service
[Unit]
Description=Corporate Memory Knowledge Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=data-ops
WorkingDirectory=/opt/data-analyst/repo
ExecStart=/usr/local/bin/collect-knowledge

# Environment
EnvironmentFile=/opt/data-analyst/.env
EnvironmentFile=/opt/data-analyst/repo/.env
Environment=CONFIG_DIR=/opt/data-analyst/instance/config
Environment=PYTHONPATH=/opt/data-analyst/repo

# Security hardening - root needed to read /home/*/CLAUDE.local.md
ProtectSystem=strict
ReadWritePaths=/data/corporate-memory
PrivateTmp=true

# Timeout for long-running HAIKU calls
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
