# server/corporate-memory.timer
[Unit]
Description=Run Corporate Memory Knowledge Collector every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
