# server/session-collector.timer
[Unit]
Description=Run Session Transcript Collector every 6 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
