# server/webapp.service
[Unit]
Description=Data Broker Web App (Google SSO)
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=webapp
WorkingDirectory=/opt/data-analyst/repo
Environment="PATH=/opt/data-analyst/.venv/bin"
EnvironmentFile=/opt/data-analyst/.env
ExecStart=/opt/data-analyst/.venv/bin/gunicorn \
    --workers 2 \
    --bind unix:/run/webapp/webapp.sock \
    --access-logfile /opt/data-analyst/logs/webapp-access.log \
    --error-logfile /opt/data-analyst/logs/webapp-error.log \
    webapp.app:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
