# services/catalog_refresh/systemd/catalog-refresh.service
[Unit]
Description=Catalog Refresh - export metrics and tables from OpenMetadata to YAML
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=data-ops
WorkingDirectory=/opt/data-analyst/repo
ExecStart=/opt/data-analyst/.venv/bin/python3 -m src.catalog_export

# Environment
EnvironmentFile=/opt/data-analyst/.env
Environment=PYTHONPATH=/opt/data-analyst/repo
Environment=CONFIG_DIR=/opt/data-analyst/instance/config

# Write access to docs output directory
ProtectSystem=strict
ReadWritePaths=/data/docs /opt/data-analyst/logs
PrivateTmp=true

# Catalog export is fast (seconds)
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
